# scripts/deploy-agents/systemd/vector-tier1.service
[Unit]
Description=Vector Module-1 agent watcher (tier-1: event-driven)
Documentation=file://%h/code/vector-agent-modules/scripts/deploy-agents/README.md
After=network-online.target
Wants=network-online.target
# Stop the restart loop after 10 failures in 1h. If this trips, something
# is systemically broken — most likely: OAuth token expired/missing/
# malformed, network to Ogmios is down, or Python deps got uninstalled.
# After StartLimitBurst failures within StartLimitIntervalSec the unit
# enters `failed` state; `systemctl --user status` shows it and the
# operator can investigate rather than burning logs forever.
StartLimitIntervalSec=3600
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=%h/vector-agents
# EnvironmentFile is loaded if present; allows the operator to set
# CLAUDE_CODE_OAUTH_TOKEN without editing this unit.
EnvironmentFile=-%h/vector-agents/env
ExecStart=/usr/bin/env python3 %h/vector-agents/bin/tier1_watcher.py
Restart=on-failure
RestartSec=30
# Let the watcher finish its current tick before systemd kills it.
TimeoutStopSec=60
KillSignal=SIGTERM
KillMode=control-group
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/vector-agents
ProtectHome=read-only
BindReadOnlyPaths=%h/code/vector-agent-modules
Umask=0077
StandardOutput=append:%h/vector-agents/logs/tier1-watcher.log
StandardError=append:%h/vector-agents/logs/tier1-watcher.log

[Install]
WantedBy=default.target
